# udevrules/99-monitor-hdmi-plugged.rules
ACTION=="change", SUBSYSTEM=="drm", RUN+="/usr/bin/fillHdmi.sh"
